# docs/systemd-runner.service
[Unit]
Description=Market Maker Runner (single symbol)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WorkingDirectory=/opt/market-maker
User=runner
Group=runner
Environment="CONFIG_PATH=/opt/market-maker/configs/round8_survival.yaml"
Environment="METRICS_ADDR=:9101"
Environment="DRY_RUN=false"
ExecStart=/opt/market-maker/bin/runner -config=${CONFIG_PATH} -metricsAddr=${METRICS_ADDR}
ExecStop=/opt/market-maker/scripts/graceful_shutdown.sh
PIDFile=/opt/market-maker/logs/runner.pid
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=45
Restart=on-failure
RestartSec=10
LimitNOFILE=65535
RuntimeDirectory=market-maker
RuntimeDirectoryMode=0755
StandardOutput=append:/var/log/market-maker/runner.log
StandardError=append:/var/log/market-maker/runner.log

[Install]
WantedBy=multi-user.target
